# crypto_dashboard/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log

[program:calculations]
command=python calculations.py
autostart=true
autorestart=true
stderr_logfile=/var/log/calculations.err.log
stdout_logfile=/var/log/calculations.out.log

[program:categoryrank]
command=python CategoryRank.py
autostart=true
autorestart=true
stderr_logfile=/var/log/categoryrank.err.log
stdout_logfile=/var/log/categoryrank.out.log

[program:mcrolling]
command=python MCrolling.py
autostart=true
autorestart=true
stderr_logfile=/var/log/mcrolling.err.log
stdout_logfile=/var/log/mcrolling.out.log

[program:strength]
command=python strength.py
autostart=true
autorestart=true
stderr_logfile=/var/log/strength.err.log
stdout_logfile=/var/log/strength.out.log

[program:tokenstrength]
command=python tokenstrength.py
autostart=true
autorestart=true
stderr_logfile=/var/log/tokenstrength.err.log
stdout_logfile=/var/log/tokenstrength.out.log